# fs/etc/services.d/pomerium/run
#!/usr/bin/with-contenv ash

exec runas /bin/pomerium $POMERIUM_ARGS